# cron/crontab
*   *   *   *   *   /usr/bin/flock -xn /var/lock/sagittariidae/cron.lock ${HOME}/sagittariidae-ws.git/cron/sweep
