# cron/crontab
*   *   *   *   *   ${HOME}/sagittariidae-ws.git/cron/sweep
